# dbackup-archive.timer
[Unit]
Description=Database Backup Archive Timer

[Timer]
OnCalendar=*-*-*/3
Persistent=true

[Install]
WantedBy=timers.target